# crontab
# Run attendance check every 5 minutes
*/5 * * * * curl -X POST http://localhost:3000/api/cron/attendance-check >> /var/log/cron/alerts.log 2>&1

# Run daily report at 11:59 PM
59 23 * * * curl -X POST http://localhost:3000/api/cron/daily-report >> /var/log/cron/alerts.log 2>&1